# infra/systemd/codex-sync.service
[Unit]
Description=Sync Codex CLI credentials from host to job-finder seed dirs

[Service]
Type=oneshot
Environment=JOB_FINDER_BOT_DIR=%h/Development/job-finder-bot
ExecStart=/usr/bin/env bash -c "${JOB_FINDER_BOT_DIR}/scripts/sync_codex_credentials.sh"
WorkingDirectory=${JOB_FINDER_BOT_DIR}

[Install]
WantedBy=default.target
