# infra/systemd/codex-sync.service
[Unit]
Description=Sync Codex CLI credentials from host to job-finder seed dirs

[Service]
Type=oneshot
ExecStart=%h/Development/job-finder-bot/scripts/sync_codex_credentials.sh
WorkingDirectory=%h/Development/job-finder-bot

[Install]
WantedBy=default.target
