# job-finder-worker/docker/crontab
# Job Finder Cron Schedule
# Run job search every 6 hours at minute 0
# Uses gosu to run as node user (uid 1000) to match codex credential ownership
0 */6 * * * root /usr/sbin/gosu node /bin/bash -c 'cd /app && /home/node/.local/bin/python scheduler.py' >> /var/log/cron.log 2>&1

# Run maintenance every 24 hours at 3 AM UTC
# - Deletes job matches older than 2 weeks
# - Recalculates match scores based on staleness (no AI)
0 3 * * * root curl -s -X POST http://localhost:5555/maintenance >> /var/log/maintenance.log 2>&1

# Run at specific times daily (optional - comment out if using above)
# 0 6 * * * root /usr/sbin/gosu node /bin/bash -c 'cd /app && /home/node/.local/bin/python scheduler.py' >> /var/log/cron.log 2>&1
# 0 12 * * * root /usr/sbin/gosu node /bin/bash -c 'cd /app && /home/node/.local/bin/python scheduler.py' >> /var/log/cron.log 2>&1
# 0 18 * * * root /usr/sbin/gosu node /bin/bash -c 'cd /app && /home/node/.local/bin/python scheduler.py' >> /var/log/cron.log 2>&1

# Blank line required at end of cron file
